# build/immortalwrt_Redmi_AX6S/files/etc/crontabs/root

#*/1 * * * * sh /etc/config/sh/ipv6ping.sh &> /dev/null
#0 */2 * * * sh /etc/config/sh/crontab.sh &> /dev/null
#0 4 * * * /etc/init.d/firewall restart && sh /etc/config/sh/bond-lan-up.sh > /dev/null
0 4 * * 1,4 sh /etc/config/sh/re.sh &> /dev/null
*/3 * * * * sh /etc/config/sh/dynv6.sh


0 2 * * * /usr/share/UnblockNeteaseMusic/update_core.sh 2>&1
